# deploy/systemd/claude-channels.service
[Unit]
Description=Claude Code with Telegram Channel (tmux-supervised)
After=network.target

[Service]
Type=forking
User=forwarder
Environment=HOME=/home/forwarder
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/forwarder/.bun/bin
ExecStartPre=-/usr/bin/tmux kill-session -t claude
ExecStart=/home/forwarder/app/deploy/run_claude_channels.sh
ExecStop=/usr/bin/tmux kill-session -t claude
# Check every 30s that the tmux session and claude process are alive.
# If either dies, systemd restarts us.
WatchdogSec=60
NotifyAccess=all
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
