# deploy/systemd/telegram-tracker.service
[Unit]
Description=Telegram Pick Grader
After=network.target

[Service]
Type=oneshot
User=forwarder
WorkingDirectory=/home/forwarder/app
EnvironmentFile=/home/forwarder/app/.env
EnvironmentFile=-/home/forwarder/app/.env.local
ExecStart=/home/forwarder/app/run_tracker.sh
SyslogIdentifier=tracker
StandardOutput=journal
StandardError=journal
