# deploy/systemd/trent-monitor.timer
[Unit]
Description=Run Trent Pick Monitor every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
AccuracySec=30s

[Install]
WantedBy=timers.target
