# deploy/systemd/trent-monitor.service
[Unit]
Description=Trent Pick Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=forwarder
Group=forwarder
WorkingDirectory=/home/forwarder/app
EnvironmentFile=/home/forwarder/app/.env
EnvironmentFile=/home/forwarder/app/.env.local
ExecStart=/home/forwarder/app/run_trent_watcher.sh
# Must exceed the runner's own worst case: attempt 1 (90s fetch timeout)
# + 60s backoff + attempt 2 (90s) = ~240s. At the old 120s systemd killed the
# retry mid-flight, so a slow/rate-limited fetch could never complete.
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
